# model_tuner.service
[Unit]
Description=Model Tuner Service

[Service]
ExecStart= /start_tuning.sh

[Install]
WantedBy=multi-user.target